# files/etc/network/if-up.d/restore-routing
#!/bin/sh
/usr/local/sbin/routes-enforce
exit 0
